# cron.daily/ntpdate
#!/bin/sh
/usr/sbin/ntpdate -s ntp.ubuntu.com
