# cron.d/softaculous
37 3 * * * root /usr/local/vesta/php/bin/php -d disable_functions="" /usr/local/vesta/softaculous/cron.php >> /dev/null 2>&1
*/5 * * * *  root /usr/local/vesta/php/bin/php -d disable_functions="" /usr/local/vesta/softaculous/do_backups.php >> /dev/null 2>&1

